# git-stats/supervisord.d/generatord.ini
[program:generatord]
command=/generatord.sh
autorestart=true
