# webproxy/supervisord.d/dockergen.ini
[program:dockergen]
command=/docker-gen -config /docker-gen.conf
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=20
