# nginx-simple/supervisord.d/inotify.ini
[program:inotify]
command=/bin/bash -c 'while inotifywait -e close_write /nginx.conf; do supervisorctl restart nginx; done'
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=10

