# ns/root/etc/services.d/inotify/run
#!/usr/bin/with-contenv sh

/usr/bin/k8-inotify-reload "/etc/dnsmasq.d/" "s6-svc -t /services/dnsmasq/"

