# sentry/supervisord.d/sentry-web.ini
[program:sentry-web]
directory=/.../
command=/.../sentry start
autostart=true
autorestart=true
redirect_stderr=true

