# sentry/supervisord.d/sentry-web.ini
[program:sentry-web]
directory=/data/
command=/usr/bin/sentry --config=/data/sentry.conf.py start
autostart=true
autorestart=true
redirect_stderr=true

