# taggo/supervisord.d/inotify.ini
[program:inotify]
command=/bin/bash -c 'while inotifywait -e move -e create -e delete -r ${TAGGO_WATCH:-/data/}; do /usr/bin/python /taggo/taggo run_once; done'
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=10
autostart=false
