# sickbeard/supervisord.d/sickbeard.ini
[program:sickbeard]
command=/usr/bin/python /sickbeard/SickBeard.py
autorestart=true
stopwaitsecs=2
stopasgroup=true
redirect_stdout=true
redirect_stderr=true
