# h5ai/supervisord.d/php-fpm.ini
[program:phpfpm]
command=/usr/sbin/php-fpm
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=20
