# buildbot/supervisord.d/buildbot-slave.ini
[program:buildbot-slave]
command=/usr/bin/buildslave start --nodaemon /data/buildbot/
stdout_logfile=/var/log/supervisor-%(program_name)s.log
stderr_logfile=/var/log/supervisor-%(program_name)s.log
autorestart=true
autostart=false
