# buildbot/supervisord.d/buildbot-master.ini
[program:buildbot-master]
command=/usr/bin/buildbot start --nodaemon /data/buildbot/
stdout_logfile=/var/log/supervisor-%(program_name)s.log
stderr_logfile=/var/log/supervisor-%(program_name)s.log
autorestart=true
autostart=false
