# nginx-proxy-front/supervisord.d/generator_inotify.ini
[program:generator_inotify]
command=/bin/bash -c 'while inotifywait -e create -e delete -e close_write -e move -r /domains/; do /bin/sh /generator.sh; done'
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=10

