# ns/supervisord.d/dnsmasq.ini
[program:dnsmasq]
command=dnsmasq -k
autorestart=true
stopwaitsecs=1
stopasgroup=true
redirect_stdout=true
redirect_stderr=true
priority=10
