# ns/supervisord.d/inotify.ini
[program:inotify]
command=/bin/bash -c 'while inotifywait -e close_write -r /etc/dnsmasq.d/; do pkill -9 dnsmasq; done'
autorestart=true
stopwaitsecs=2
stopasgroup=true
redirect_stdout=true
redirect_stderr=true
priority=20
