# wetty/supervisord.d/wetty.ini
[program:wetty]
command=/usr/bin/node /wetty/app.js -p 3000
autorestart=true
stopwaitsecs=2
stopasgroup=true
