# etherpad/supervisord.d/etherpad.ini
[program:etherpad]
command=/usr/bin/node node_modules/ep_etherpad-lite/node/server.js
directory=/etherpad-lite
autorestart=true
stopwaitsecs=1
stopasgroup=true
redirect_stdout=true
redirect_stderr=true
