# flexget/supervisord.d/flexget.ini
[program:flexget]
command=/usr/bin/flexget -c /config.yml daemon start
autorestart=true
stopwaitsecs=2
stopasgroup=true
redirect_stdout=true
redirect_stderr=true
