# zabbix/supervisord.d/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -e info -DFOREGROUND
autorestart=true
stopwaitsecs=2
