# noted/supervisord.d/ipython_notebook.ini
[program:ipython_notebook]
command=/usr/bin/scl enable python33 'ipython notebook --no-browser --profile=nbserver'
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=20
