# noted/supervisord.d/taggo_timer.ini
[program:taggo_timer]
command=/bin/bash -c 'while true; do /usr/bin/python /taggo/taggo run_once; sleep ${TAGGO_INTERVAL:-60}; done'
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=10
