# noted/supervisord.d/taggo.ini
[program:taggo]
command=
autorestart=true
stopwaitsecs=2
stopasgroup=true
priority=10
