# centos-base/supervisord.d/appname.ini.sample
[program:sample]
command=/bin/bash -c 'while true; do echo "running.."; sleep 1; done'
autorestart=true
stopwaitsecs=2
stopasgroup=true
redirect_stdout=true
redirect_stderr=true
